# cleanup.timer
[Unit]
Description=Daily Webhook Cleanup
Requires=cleanup.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target